# runit/amiws/run
#!/bin/sh

exec 2>&1

mkdir -p /run/amiws

config=${AMIWS_CONFIG:-/data/amiws/config.yaml}

exec amiws -f ${config}
